# data/claude-logs/com.athena.claude-logs-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.athena.claude-logs-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/j/Code/athena/data/claude-logs/sync_claude_logs.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>300</integer> <!-- Run every 5 minutes (300 seconds) -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/j/Code/athena/data/claude-logs/sync.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/j/Code/athena/data/claude-logs/sync.stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/j/Code/athena/data/claude-logs</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>